# deploy/ai.plutus.updater.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- Template: scripts/install_launchd.sh substitutes __ROOT__ / __LOGDIR__. -->
<plist version="1.0">
<dict>
  <key>Label</key><string>ai.plutus.updater</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>__ROOT__/scripts/update.sh</string>
  </array>
  <key>WorkingDirectory</key><string>__ROOT__</string>
  <key>StartInterval</key><integer>21600</integer>
  <key>ProcessType</key><string>Background</string>
  <key>StandardOutPath</key><string>__LOGDIR__/update.log</string>
  <key>StandardErrorPath</key><string>__LOGDIR__/update.err.log</string>
</dict>
</plist>
